# whitepine/crontab
# Runner
MAILTO=""
03 7 * * * "$HOME/bin/runner.sh"
